# templates/mesh_bot.service
[Unit]
Description=Meshing-Around Mesh Bot (mesh_bot.py)
Documentation=https://github.com/SpudGunMan/meshing-around
After=network.target meshtasticd.service
Wants=meshtasticd.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__BOT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash __BOT_DIR__/launch.sh mesh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
